# ominfra/supervisor/tests/configs/demo.conf
[supervisord]
no_daemon=true
# log_file=/dev/null
# log_file_max_bytes=4194304

[program:sleep60]
command=sleep 60
stdout_log_file=/dev/fd/1
stdout_log_file_max_bytes=0
redirect_stderr=true

[program:sleep600]
command=sleep 600
stdout_log_file=/dev/fd/1
stdout_log_file_max_bytes=0
redirect_stderr=true
process_name=sleep600.%(process_num)d
num_procs=2

[group:sleeps]
programs=sleep60,sleep600

[program:ls]
command=ls -al
stdout_log_file=/dev/fd/1
stdout_log_file_max_bytes=0
redirect_stderr=true
